# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_error.log
priority=1

; Replaces the old three Python processes (config-ui, fdroid-updater,
; patch-checker) - the admin server and both scheduler loops (GitHub
; releases, APKMirror patching) now run as coroutines inside this one JVM
; process (see morphe-fdroid-server/app/src/main/kotlin/app/fdroidserver/Main.kt).
; morphe-patcher is a normal Gradle dependency merged into app.jar's own
; shadow/fat jar (see morphe-fdroid-server/app/build.gradle.kts) - no second
; jar on the classpath.
[program:fdroid-server]
command=java -jar /app/app.jar
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/fdroid-server.log
stderr_logfile=/var/log/supervisor/fdroid-server_error.log
priority=2
